# deployment/storage-bucket-frontend.service
[Unit]
Description=Storage Bucket Frontend
After=network.target storage-bucket-backend.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/storage-bucket/frontend
Environment="NODE_ENV=production"
Environment="PORT=3200"
EnvironmentFile=/var/www/storage-bucket/frontend/.env
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=storage-bucket-frontend

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
