# deployment/storage-bucket-backend.service
[Unit]
Description=Storage Bucket Backend API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/storage-bucket/backend
Environment="NODE_ENV=production"
Environment="PORT=3201"
EnvironmentFile=/var/www/storage-bucket/backend/.env
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=storage-bucket-backend

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
